# myapp.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/ec2-user/coffee-shop/startserver

[Install]
WantedBy=default.target